# files/etc/systemd/system/ecos-power-profile.service
[Unit]
Description=ecos — Apply gaming power profile
After=tuned.service

[Service]
Type=oneshot
ExecStart=/bin/tuned-adm profile ecos-gaming

[Install]
WantedBy=multi-user.target