# files/etc/systemd/system/ecos-firstboot.service
[Unit]
Description=ecOS — First boot configuration
After=local-fs.target
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ecos-firstboot.sh
RemainAfterExit=yes
# Ne pas bloquer le boot si le script échoue
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target